# spec/debian/twindb-server-dispatcher.cron.d
@hourly                     root    twindb-apply_retention_policy.php
@hourly                     root    twindb-volume_usage_snapshot.php
*/5 * * * *                 root    twindb-schedule.php